# ha-addon/rootfs/etc/services.d/esphome-dist-server/finish
#!/usr/bin/execlineb -S1
/run s6-svscanctl -t /var/run/s6/services
