# ha-addon/rootfs/etc/services.d/esphome-dist-server/run
#!/bin/sh
echo "Starting Fleet for ESPHome"
exec python3 /app/main.py
